# rootbase/etc/inittab
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -a
::sysinit:/etc/init.d/rcS

console::respawn:/sbin/getty -L console 0 vt100
tty2::respawn:/sbin/getty 38400 tty2
tty3::respawn:/sbin/getty 38400 tty3
tty4::respawn:/sbin/getty 38400 tty4

::shutdown:/etc/init.d/rcK
::shutdown:/bin/umount -a -r
